# autostart/skiller_backslash.desktop
[Desktop Entry]
Name=Backslash is backslash instead of insert
Exec=xmodmap -e "keycode 118 = backslash bar"
Terminal=false
Type=Application
