# autostart/skiller_insert.desktop
[Desktop Entry]
Name=Insert becomes backspace
Exec=xmodmap -e "keycode 51 = BackSpace"
Terminal=false
Type=Application